# stage1/rootfs/aggregate/units/ctrl-alt-del.target
[Unit]
Description=Rocket SIGINT target
DefaultDependencies=false
Requires=sigint.service
