# stage1/rootfs/aggregate/units/sigint.service
[Unit]
Description=Rocket SIGINT isolator
DefaultDependencies=false

[Service]
ExecStart=/usr/bin/systemctl isolate exit-watcher.service
